# overlay2/6b28a0a1a022b710a9bae6ef38784ca01fdc879e91ab5499b096bb52703faed4/merged/lib/systemd/system/apache-htcacheclean.service
[Unit]
Description=Disk Cache Cleaning Daemon for Apache HTTP Server
After=apache2.service
Documentation=https://httpd.apache.org/docs/2.4/programs/htcacheclean.html

[Service]
Type=forking
User=www-data
Environment=HTCACHECLEAN_SIZE=300M
Environment=HTCACHECLEAN_DAEMON_INTERVAL=120
Environment=HTCACHECLEAN_PATH=/var/cache/apache2/mod_cache_disk
Environment=HTCACHECLEAN_OPTIONS=-n
EnvironmentFile=-/etc/default/apache-htcacheclean
ExecStart=/usr/bin/htcacheclean -d $HTCACHECLEAN_DAEMON_INTERVAL -p $HTCACHECLEAN_PATH -l $HTCACHECLEAN_SIZE $HTCACHECLEAN_OPTIONS

[Install]
WantedBy=multi-user.target
